# OpenBSD/rc.d/crsh
#!/bin/ksh
#
# $OpenBSD: crsh-helper,v 1.2 2020/02/20 19:52:12 olevole Exp $

daemon="/usr/bin/crsh-helper.sh"

. /etc/rc.d/rc.subr

rc_cmd $1
